# etc/supervisord.d/influxdb.ini
[program:influxdb]
command=/opt/qnib/bin/start_influxdb.sh
redirect_stderr=true
stdout_logfile=syslog
autorestart = true
priority = 10
