# etc/supervisord.d/influxdb-init.ini
[program:influxdb-init]
command=/opt/qnib/influxdb/bin/init.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
autorestart = false
priority = 10
