# etc/supervisord.d/influxdb.ini
[program:influxdb]
command=/opt/qnib/influxdb/bin/start.sh
redirect_stderr=true
#stdout_logfile=syslog
stdout_logfile=/var/log/supervisor/%(program_name)s.log
autorestart = true
priority = 10
