# etc/supervisord.d/influxdb.ini
[program:influxdb]
command=/usr/bin/influxdb -pidfile /opt/influxdb/shared/influxdb.pid -config /opt/influxdb/shared/config.toml
redirect_stderr=true
stdout_syslog=true
autorestart = true
priority = 10
